# gnu/lib/dbus/bus/dbus.socket.in
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=@DBUS_SYSTEM_SOCKET@
