# overlays/uzip/localize/files/usr/local/etc/rc.d/localize
#!/bin/sh

# PROVIDE: localize
# REQUIRE: DAEMON NETWORKING
# BEFORE:  LOGIN
# KEYWORD: nojail shutdown

PATH=$PATH:/usr/local/sbin
export PATH

. /etc/rc.subr

name="localize"
start_cmd="${name}_start"

localize_start()
{
    /usr/local/sbin/localize
}

load_rc_config $name
run_rc_command "$1"
